# config/claude-telegram.service
[Unit]
Description=Claude Code with Telegram Channel
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
User=REPLACE_USER
Environment=PATH=/home/REPLACE_USER/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/REPLACE_USER
WorkingDirectory=/home/REPLACE_USER
ExecStartPre=-/usr/bin/tmux kill-session -t claude-telegram
ExecStart=/usr/bin/tmux new-session -d -s claude-telegram '/usr/local/bin/claude --channels plugin:telegram@claude-plugins-official --permission-mode acceptEdits'
ExecStop=/usr/bin/tmux kill-session -t claude-telegram
RemainAfterExit=yes
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
